# RaspberryPi/VMCP/loopplay.conf
[program:loopplay]
command=/home/pi/loopplay.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/loopplay.err.log
#stdout_logfile=/var/log/loopplay.out.log

