# RaspberryPi/VMCP/monitor.conf
[program:monitor]
command=/home/pi/monitor.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/monitor.err.log
stdout_logfile=/var/log/monitor.out.log

